# sdk/rockx-x86-64-Linux/bin/npu_transfer_proxy/npu_transfer_proxy.service
[Unit]
Description=NPU Transfer Proxy for Rockchip RK1808 AI Compute Stick

[Service]
Restart=always
ExecStart=/usr/local/bin/npu_transfer_proxy

[Install]
WantedBy=default.target
